# conf/app_supervisor.conf
[program:flaskapp]
command = uwsgi --ini /conf/app_uwsgi.ini

[program:nginx]
command = service nginx restart

[program:redis]
command = service redis-server restart
